# powerai-examples/conf.d/tensorflow_jupyter.conf
[program:tensorflow]
command=/root/startjupyter.sh 8890
directory=/opt/DL/tensorflow/
redirect_stderr=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
auto_start=true
autorestart=true

