# powerai-examples/conf.d/tensorboard.conf
[program:tensorboard]
command=/root/starttensorboard.sh
directory=/opt/DL/tensorflow/
redirect_stderr=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
auto_start=true
autorestart=true
killasgroup=true

